# resources/mysupervisor/conf.d/octane.ini
[program:laravel-octane]
command=php /var/www/artisan octane:start --host=0.0.0.0 --port=80
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0