# resources/mysupervisor/conf.d/schedule.ini
[program:laravel-schedule]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0